# services/monitor_v4/run
#!/bin/bash
. $CBROOTLOCAL/setup.env
sleep 1
exec </dev/null
exec 2>&1
echo Starting
exec init_perlbrew ~/ntppool/monitor $MONITOR_API_V4 $MONITOR_API_URL
